# docs/docker/corejj.ini
[supervisord]
nodaemon=true

[program:nginx]
command = nginx -g "daemon off;"
autostart = true

[program:redis]
command = redis-server /etc/redis.conf
autostart = true

[program:corejj-api]
command = sh /data/corejj/corejj_api/tools/start-api.sh
autostart = true
stdout_logfile = /data/corejj/corejj_api/logs/api.log
redirect_stderr = true

[program:corejj-ws]
command = sh /data/corejj/corejj_api/tools/start-ws.sh
autostart = true
stdout_logfile = /data/corejj/corejj_api/logs/ws.log
redirect_stderr = true

[program:corejj-worker]
command = sh /data/corejj/corejj_api/tools/start-worker.sh
autostart = true
stdout_logfile = /data/corejj/corejj_api/logs/worker.log
redirect_stderr = true

[program:corejj-monitor]
command = sh /data/corejj/corejj_api/tools/start-monitor.sh
autostart = true
startsecs = 3
stdout_logfile = /data/corejj/corejj_api/logs/monitor.log
redirect_stderr = true

[program:corejj-scheduler]
command = sh /data/corejj/corejj_api/tools/start-scheduler.sh
autostart = true
startsecs = 3
stdout_logfile = /data/corejj/corejj_api/logs/scheduler.log
redirect_stderr = true
