# corejj_api/tools/supervisor-corejj.ini
[program:corejj-api]
command = bash /data/corejj-v3.3.3/corejj_api/tools/start-api.sh
autostart = true
stdout_logfile = /data/corejj-v3.3.3/corejj_api/logs/api.log
redirect_stderr = true

[program:corejj-ws]
command = bash /data/corejj-v3.3.3/corejj_api/tools/start-ws.sh
autostart = true
stdout_logfile = /data/corejj-v3.3.3/corejj_api/logs/ws.log
redirect_stderr = true

[program:corejj-worker]
command = bash /data/corejj-v3.3.3/corejj_api/tools/start-worker.sh
autostart = true
stdout_logfile = /data/corejj-v3.3.3/corejj_api/logs/worker.log
redirect_stderr = true

[program:corejj-monitor]
command = bash /data/corejj-v3.3.3/corejj_api/tools/start-monitor.sh
autostart = true
stdout_logfile = /data/corejj-v3.3.3/corejj_api/logs/monitor.log
redirect_stderr = true

[program:corejj-scheduler]
command = bash /data/corejj-v3.3.3/corejj_api/tools/start-scheduler.sh
autostart = true
stdout_logfile = /data/corejj-v3.3.3/corejj_api/logs/scheduler.log
redirect_stderr = true
